# systemd/update-checker.timer
[Unit]
Description=Daily Arch update check at 09:00

[Timer]
OnCalendar=09:00
Persistent=true

[Install]
WantedBy=timers.target
