# systemd/update-checker.service
[Unit]
Description=Check for Arch/AUR updates and notify via Telegram
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/update-bot.env
ExecStart=%h/.local/bin/uv run %h/scripts/update_checker.py
